# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

; ── FastAPI backend (port 8000) ──────────────────────────────────────────
[program:api]
command=/opt/venv/bin/uvicorn api:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/api.log
stderr_logfile=/var/log/supervisor/api.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; ── Streamlit frontend (port 8501) ───────────────────────────────────────
[program:streamlit]
command=/opt/venv/bin/streamlit run app.py
    --server.port=8501
    --server.address=0.0.0.0
    --server.headless=true
    --server.enableCORS=false
environment=BACKEND_URL="http://localhost:8000"
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/streamlit.log
stderr_logfile=/var/log/supervisor/streamlit.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB