# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:backend]
command=uvicorn api:app --host 0.0.0.0 --port 8000
autostart=true
autorestart=true
stderr_logfile=/var/log/backend.err.log
stdout_logfile=/var/log/backend.out.log

[program:frontend]
# We map Streamlit to the PORT environment variable that Render injects
command=sh -c 'streamlit run app.py --server.port ${PORT:-7860} --server.address 0.0.0.0'
autostart=true
autorestart=true
stderr_logfile=/var/log/frontend.err.log
stdout_logfile=/var/log/frontend.out.log